# debian/fake-hwclock.cron.hourly
#!/bin/sh
#
# Simple cron script - save the current clock periodically in case of
# a power failure or other crash

/bin/fake-hwclock save
